# system/backup-scripts/cron.daily/srvbackup
#!/bin/bash
#

/srv/scripts/backup-all.sh
